# .config/autostart/redshift-gtk.desktop
[Desktop Entry]
Exec=redshift-gtk
Comment=Color temperature adjustment tool
Version=1.0
Name=Redshift
Name[zh_CN]=红移
Name[fr]=Redshift
Name[de]=Redshift
GenericName[zh_CN]=色温调节
Name[hu]=Redshift
GenericName[tr]=Renk sıcaklığı ayarı
Comment[cs]=Nástroj pro přizpůsobení barevné teploty
Name[pl]=Redshift
Terminal=false
GenericName[da]=Justering af farvetemperatur
GenericName[it]=Regolazione della temperatura del colore
GenericName[fr]=Réglage de la température de couleur
Name[tr]=Redshift
GenericName[de]=Farbtemperatur Anpassung
Comment[hu]=Színhőmérséklet beállító eszköz
Comment[nl]=Kleurtemperatuur aanpassingsprogramma
Comment[tr]=Renk sıcaklığı ayarlama aracı
Comment[zh_CN]=色温调节工具
Comment[lt]=Spalvos temperatūros reguliavimo įrankis
GenericName[ru]=Настройка цветовой температуры
Comment[pl]=Narzędzie do dostosowywania temperatury barwowej
GenericName=Color temperature adjustment
Comment[de]=Farbtemperatur Anpassungswerkzeug
Icon=redshift
Name[it]=Redshift
GenericName[lt]=Spalvos temperatūros reguliavimas
Hidden=false
Name[ru]=Redshift
X-GNOME-Autostart-enabled=true
GenericName[nl]=Kleurtemperatuur aanpassing
GenericName[cs]=Přizpůsobení barevné teploty
Comment[da]=Justeringsværktøj for farvetemperatur
Name[es]=Redshift
Name[pt_BR]=Redshift
GenericName[hu]=Szín hőmérséklet beállítás
GenericName[pl]=Dostosowanie temperatury barwowej
Name[lt]=Redshift
StartupNotify=true
Comment[es]=Herramienta para el ajuste de la temperatura de color
Comment[it]=Strumento per la regolazione della temperatura del colore
GenericName[es]=Ajuste de la temperatura de color
Comment[fr]=Outil de réglage de la température de couleur
Comment[ru]=Инструмент регулирования цветовой температуры
Type=Application
Categories=Utility;

